# infra/marina-bridge/cmd/emtrak-forwarder/systemd/emtrak-forwarder.service
[Unit]
Description=em-trak AIS forwarder (NMEA 0183 over WiFi -> ais-cache)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/emtrak-forwarder.env
ExecStart=/usr/local/bin/emtrak-forwarder
Restart=always
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
